# etc/cron.d/cron-smeserver
45 */4 * * * root rsync --progress -aHS --partial --delete-after smeserver.bhs.mirrors.ovh.net::smeserver /home/e-smith/files/ibays/smeserver/html/ >/dev/null 2>&1 
