# etc/cron.d/sauvMysql2stephdl.dyndns.org
35 15 * * * root rsync -avz --delete /root/backup/ root@stephdl.dyndns.org:/root/mysql_backup_mirror.de-labrusse.fr
