# etc/cron.d/TarBackupDokuwiki
* 5 * * *  root /sbin/e-smith/TarBackupDokuwiki
